# overlays/beaglebone/rootfs-overlay/etc/profile.d/env.sh
#!/bin/bash
ENFORCE_VERIFY=1  # 0=warn only, 1=enforce (reboot on failure)